# rios-worker/rios-worker.service
[Unit]
Description=RiOS Worker - Decentralized GPU Compute Node
Documentation=https://docs.rios.com.ai
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=%i
Group=docker

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Worker configuration is loaded from ~/.rios/config.json
ExecStart=/usr/local/bin/rios-worker run

# Restart policy
Restart=always
RestartSec=10

# Resource limits (adjust as needed)
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rios-worker

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

